# macos-launchd/Library/LaunchAgents/com.user.killchrome.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN"
 "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.user.killchrome</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/YOUR_USERNAME/.local/bin/kill_chrome_if_running_too_long.sh</string>
  </array>

  <key>StartInterval</key>
  <integer>600</integer> <!-- every 10 minutes -->

  <key>RunAtLoad</key>
  <true/>
</dict>
</plist>
